# apache2/apache-supervisor.conf
[program:apache]
command=apachectl -X
redirect_stderr=true